# deploy/guild-hall.service
[Unit]
Description=Guild Hall - Claude Agent SDK Session Manager
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/rjroy/Projects/guild-hall
ExecStart=/usr/bin/bun run start
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=guild-hall

[Install]
WantedBy=default.target
